# rootfs/etc/s6-overlay/s6-rc.d/minidlna/up
#!/command/execlineb -P
foreground { echo "[s6-init] Iniciando MiniDLNA..." }
/etc/init.d/minidlna start
