# rootfs/etc/s6-overlay/s6-rc.d/nfs-server/up
#!/command/execlineb -P
foreground { echo "[s6-init] Starting NFS Server (NFSv3 and NFSv4)..." }
foreground { /bin/sh -c "
    # Mount nfsd if not already mounted
    if ! mountpoint -q /proc/fs/nfsd; then
        mount -t nfsd nfsd /proc/fs/nfsd 2>/dev/null || echo 'WARN: Could not mount /proc/fs/nfsd'
    fi
    
    # Check that /etc/exports exists
    if [ ! -f /etc/exports ]; then
        echo 'ERROR: /etc/exports does not exist'
        exit 1
    fi
    
    # Enable NFSv3 and NFSv4 explicitly
    # Create configuration if it does not exist
    if [ ! -f /etc/default/nfs-kernel-server ]; then
        echo 'RPCNFSDCOUNT=8' > /etc/default/nfs-kernel-server
        echo 'RPCNFSDPRIORITY=0' >> /etc/default/nfs-kernel-server
        echo 'RPCMOUNTDOPTS=\"--manage-gids\"' >> /etc/default/nfs-kernel-server
        echo 'NEED_SVCGSSD=no' >> /etc/default/nfs-kernel-server
        echo 'RPCSVCGSSDOPTS=\"\"' >> /etc/default/nfs-kernel-server
    fi
    
    # Ensure /etc/exports.d exists
    mkdir -p /etc/exports.d
    
    # Start NFS services (supports v3 and v4 by default on Debian)
    /etc/init.d/nfs-kernel-server start
    
    # Verify supported versions
    echo 'Supported NFS versions:'
    cat /proc/fs/nfsd/versions 2>/dev/null || echo 'Could not verify versions'
" }

