# rootfs/etc/s6-overlay/s6-rc.d/wsdd/up
#!/command/execlineb -P
foreground { echo "[s6-init] Iniciando WSDD..." }
/etc/init.d/wsdd start
