# rootfs/etc/s6-overlay/s6-rc.d/rpcbind/up
#!/command/execlineb -P
foreground { echo "[s6-init] Starting RPC Bind..." }
if { /etc/init.d/rpcbind start }
